# udev/rules.d/70-snap.snapd-desktop-integration.rules
# This file is automatically generated.
# opengl
KERNEL=="dma_buf_te", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="galcore", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="mali[0-9]*", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="nvhost-*", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="nvmap", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="pvr_sync", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="vchiq", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="linux,cma", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="system", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_snapd-desktop-integration_snapd-desktop-integration"
TAG=="snap_snapd-desktop-integration_snapd-desktop-integration", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_snapd-desktop-integration_snapd-desktop-integration $devpath $major:$minor"
