# udev/rules.d/70-snap.firefox.rules
# This file is automatically generated.
# camera
KERNEL=="vchiq", TAG+="snap_firefox_firefox"
# camera
KERNEL=="video[0-9]*", TAG+="snap_firefox_firefox"
# joystick
KERNEL=="event[0-9]*", SUBSYSTEM=="input", ENV{ID_INPUT_JOYSTICK}=="1", TAG+="snap_firefox_firefox"
# joystick
KERNEL=="full", SUBSYSTEM=="mem", TAG+="snap_firefox_firefox"
# joystick
KERNEL=="js[0-9]*", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="dma_buf_te", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="galcore", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="mali[0-9]*", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="nvhost-*", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="nvmap", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="pvr_sync", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="vchiq", TAG+="snap_firefox_firefox"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_firefox_firefox"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="linux,cma", TAG+="snap_firefox_firefox"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="system", TAG+="snap_firefox_firefox"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_firefox_firefox"
# u2f-devices
# Bluink Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="2abe", ATTRS{idProduct}=="1002", TAG+="snap_firefox_firefox"
# u2f-devices
# Feitian ePass FIDO, BioPass FIDO2
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="096e", ATTRS{idProduct}=="0850|0852|0853|0854|0856|0858|085a|085b|085d", TAG+="snap_firefox_firefox"
# u2f-devices
# GoTrust Idem Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="32a3", ATTRS{idProduct}=="3201|3203", TAG+="snap_firefox_firefox"
# u2f-devices
# Google Titan U2F
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="18d1", ATTRS{idProduct}=="5026|9470", TAG+="snap_firefox_firefox"
# u2f-devices
# Happlink (formerly Plug-Up) Security KEY
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="2581", ATTRS{idProduct}=="f1d0", TAG+="snap_firefox_firefox"
# u2f-devices
# HyperSecu HyperFIDO
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="096e|2ccf", ATTRS{idProduct}=="0880", TAG+="snap_firefox_firefox"
# u2f-devices
# HyperSecu HyperFIDO Pro
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="2ccf", ATTRS{idProduct}=="0854", TAG+="snap_firefox_firefox"
# u2f-devices
# JaCarta U2F
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="24dc", ATTRS{idProduct}=="0101|0501", TAG+="snap_firefox_firefox"
# u2f-devices
# Kensington VeriMark Guard Fingerprint Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="047d", ATTRS{idProduct}=="8055", TAG+="snap_firefox_firefox"
# u2f-devices
# Ledger Blue + Nano S + Nano X + Nano S+ + Ledger Stax
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="2c97", ATTRS{idProduct}=="0000|0001|0004|0005|0015|1005|1015|4005|4015|5005|5015|6005|6015", TAG+="snap_firefox_firefox"
# u2f-devices
# MIRKey
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a2ac", TAG+="snap_firefox_firefox"
# u2f-devices
# Neowave Keydo and Keydo AES
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1e0d", ATTRS{idProduct}=="f1d0|f1ae", TAG+="snap_firefox_firefox"
# u2f-devices
# Nitrokey 3
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="20a0", ATTRS{idProduct}=="42b2", TAG+="snap_firefox_firefox"
# u2f-devices
# Nitrokey FIDO U2F
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="20a0", ATTRS{idProduct}=="4287", TAG+="snap_firefox_firefox"
# u2f-devices
# Nitrokey FIDO2
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="20a0", ATTRS{idProduct}=="42b1", TAG+="snap_firefox_firefox"
# u2f-devices
# OnlyKey
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1d50", ATTRS{idProduct}=="60fc", TAG+="snap_firefox_firefox"
# u2f-devices
# RSA DS100
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="15e1", ATTRS{idProduct}=="2019", TAG+="snap_firefox_firefox"
# u2f-devices
# SoloKeys
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="5070|50b0|beee", TAG+="snap_firefox_firefox"
# u2f-devices
# Swissbit iShield Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1370", ATTRS{idProduct}=="0911", TAG+="snap_firefox_firefox"
# u2f-devices
# Thetis Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1ea8", ATTRS{idProduct}=="f025", TAG+="snap_firefox_firefox"
# u2f-devices
# Thetis U2F BT Fido2 Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1ea8", ATTRS{idProduct}=="fc25", TAG+="snap_firefox_firefox"
# u2f-devices
# Token2 FIDO2 key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="349e", ATTRS{idProduct}=="0010|0011|0012|0020|0021|0022|0200|0201|0202", TAG+="snap_firefox_firefox"
# u2f-devices
# Tomu board + chopstx U2F + SoloKeys + Flipper zero
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="cdab|a2ca|5741", TAG+="snap_firefox_firefox"
# u2f-devices
# Trezor
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="534c", ATTRS{idProduct}=="0001|0002", TAG+="snap_firefox_firefox"
# u2f-devices
# Trezor v2
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="53c0|53c1", TAG+="snap_firefox_firefox"
# u2f-devices
# TrustKey TrustKey G310H
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="311f", ATTRS{idProduct}=="4a2a", TAG+="snap_firefox_firefox"
# u2f-devices
# U2F Zero
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="10c4", ATTRS{idProduct}=="8acf", TAG+="snap_firefox_firefox"
# u2f-devices
# U2F-TOKEN (Tomu et al.)
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0e90", TAG+="snap_firefox_firefox"
# u2f-devices
# VASCO SeccureClick
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1a44", ATTRS{idProduct}=="00bb", TAG+="snap_firefox_firefox"
# u2f-devices
# Yubico YubiKey
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1050", ATTRS{idProduct}=="0113|0114|0115|0116|0120|0121|0200|0402|0403|0406|0407|0410", TAG+="snap_firefox_firefox"
TAG=="snap_firefox_firefox", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_firefox_firefox $devpath $major:$minor"
# camera
KERNEL=="vchiq", TAG+="snap_firefox_geckodriver"
# camera
KERNEL=="video[0-9]*", TAG+="snap_firefox_geckodriver"
# joystick
KERNEL=="event[0-9]*", SUBSYSTEM=="input", ENV{ID_INPUT_JOYSTICK}=="1", TAG+="snap_firefox_geckodriver"
# joystick
KERNEL=="full", SUBSYSTEM=="mem", TAG+="snap_firefox_geckodriver"
# joystick
KERNEL=="js[0-9]*", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="dma_buf_te", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="galcore", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="mali[0-9]*", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="nvhost-*", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="nvmap", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="pvr_sync", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="vchiq", TAG+="snap_firefox_geckodriver"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_firefox_geckodriver"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="linux,cma", TAG+="snap_firefox_geckodriver"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="system", TAG+="snap_firefox_geckodriver"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Bluink Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="2abe", ATTRS{idProduct}=="1002", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Feitian ePass FIDO, BioPass FIDO2
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="096e", ATTRS{idProduct}=="0850|0852|0853|0854|0856|0858|085a|085b|085d", TAG+="snap_firefox_geckodriver"
# u2f-devices
# GoTrust Idem Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="32a3", ATTRS{idProduct}=="3201|3203", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Google Titan U2F
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="18d1", ATTRS{idProduct}=="5026|9470", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Happlink (formerly Plug-Up) Security KEY
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="2581", ATTRS{idProduct}=="f1d0", TAG+="snap_firefox_geckodriver"
# u2f-devices
# HyperSecu HyperFIDO
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="096e|2ccf", ATTRS{idProduct}=="0880", TAG+="snap_firefox_geckodriver"
# u2f-devices
# HyperSecu HyperFIDO Pro
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="2ccf", ATTRS{idProduct}=="0854", TAG+="snap_firefox_geckodriver"
# u2f-devices
# JaCarta U2F
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="24dc", ATTRS{idProduct}=="0101|0501", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Kensington VeriMark Guard Fingerprint Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="047d", ATTRS{idProduct}=="8055", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Ledger Blue + Nano S + Nano X + Nano S+ + Ledger Stax
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="2c97", ATTRS{idProduct}=="0000|0001|0004|0005|0015|1005|1015|4005|4015|5005|5015|6005|6015", TAG+="snap_firefox_geckodriver"
# u2f-devices
# MIRKey
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a2ac", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Neowave Keydo and Keydo AES
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1e0d", ATTRS{idProduct}=="f1d0|f1ae", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Nitrokey 3
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="20a0", ATTRS{idProduct}=="42b2", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Nitrokey FIDO U2F
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="20a0", ATTRS{idProduct}=="4287", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Nitrokey FIDO2
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="20a0", ATTRS{idProduct}=="42b1", TAG+="snap_firefox_geckodriver"
# u2f-devices
# OnlyKey
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1d50", ATTRS{idProduct}=="60fc", TAG+="snap_firefox_geckodriver"
# u2f-devices
# RSA DS100
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="15e1", ATTRS{idProduct}=="2019", TAG+="snap_firefox_geckodriver"
# u2f-devices
# SoloKeys
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="5070|50b0|beee", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Swissbit iShield Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1370", ATTRS{idProduct}=="0911", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Thetis Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1ea8", ATTRS{idProduct}=="f025", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Thetis U2F BT Fido2 Key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1ea8", ATTRS{idProduct}=="fc25", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Token2 FIDO2 key
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="349e", ATTRS{idProduct}=="0010|0011|0012|0020|0021|0022|0200|0201|0202", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Tomu board + chopstx U2F + SoloKeys + Flipper zero
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="cdab|a2ca|5741", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Trezor
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="534c", ATTRS{idProduct}=="0001|0002", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Trezor v2
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="53c0|53c1", TAG+="snap_firefox_geckodriver"
# u2f-devices
# TrustKey TrustKey G310H
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="311f", ATTRS{idProduct}=="4a2a", TAG+="snap_firefox_geckodriver"
# u2f-devices
# U2F Zero
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="10c4", ATTRS{idProduct}=="8acf", TAG+="snap_firefox_geckodriver"
# u2f-devices
# U2F-TOKEN (Tomu et al.)
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0e90", TAG+="snap_firefox_geckodriver"
# u2f-devices
# VASCO SeccureClick
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1a44", ATTRS{idProduct}=="00bb", TAG+="snap_firefox_geckodriver"
# u2f-devices
# Yubico YubiKey
SUBSYSTEM=="hidraw", KERNEL=="hidraw*", ATTRS{idVendor}=="1050", ATTRS{idProduct}=="0113|0114|0115|0116|0120|0121|0200|0402|0403|0406|0407|0410", TAG+="snap_firefox_geckodriver"
TAG=="snap_firefox_geckodriver", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_firefox_geckodriver $devpath $major:$minor"
