# systemd/system/snap-snapd\x2ddesktop\x2dintegration-178.mount
[Unit]
Description=Mount unit for snapd-desktop-integration, revision 178
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/snapd-desktop-integration_178.snap
Where=/snap/snapd-desktop-integration/178
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
