# systemd/system/snap-snapd-20671.mount
[Unit]
Description=Mount unit for snapd, revision 20671
After=snapd.mounts-pre.target
Before=snapd.mounts.target
Before=local-fs.target

[Mount]
What=/var/lib/snapd/snaps/snapd_20671.snap
Where=/snap/snapd/20671
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
