# systemd/system/snap-snapd-21759.mount
[Unit]
Description=Mount unit for snapd, revision 21759
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/snapd_21759.snap
Where=/snap/snapd/21759
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
